# templates/redis-master.service.j2
[Unit]
Description=Redis Master
After=network.target

[Service]
Type=simple
PIDFile={{ redis_master_run_dir }}/redis.pid
User=root
ExecStartPre=/bin/mkdir -p {{ redis_master_run_dir }}/logs
ExecStart=/data/redis/master/src/redis-server {{ redis_master_run_dir }}/redis.conf
ExecStop=/bin/kill -15 $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
